# SessionManager/launchd/ai.claude.symlink-guard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>ai.claude.symlink-guard</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/alex/A/SessionManager/scripts/ensure-claude-symlinks.sh</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>StartInterval</key><integer>3600</integer>
  <key>StandardOutPath</key><string>/Users/alex/Library/Logs/claude-symlink-guard.launchd.log</string>
  <key>StandardErrorPath</key><string>/Users/alex/Library/Logs/claude-symlink-guard.launchd.log</string>
</dict>
</plist>
