# SessionManager/launchd/ai.sessionskills.ingest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>ai.sessionskills.ingest</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/alex/A/SessionManager/bin/sm-pipeline</string>
    <string>run</string><string>ingest</string><string>merge</string>
  </array>
  <key>StartInterval</key><integer>900</integer>
  <key>RunAtLoad</key><false/>
  <key>StandardOutPath</key><string>/Users/alex/Library/Logs/sessionskills-ingest.log</string>
  <key>StandardErrorPath</key><string>/Users/alex/Library/Logs/sessionskills-ingest.log</string>
  <key>WorkingDirectory</key><string>/Users/alex/A/SessionManager</string>
</dict>
</plist>
