# SessionManager/launchd/ai.sessionskills.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>ai.sessionskills.weekly</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/alex/A/SessionManager/bin/sm-pipeline</string>
    <string>run</string><string>prune</string><string>archive</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>0</integer>
    <key>Hour</key><integer>3</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>RunAtLoad</key><false/>
  <key>StandardOutPath</key><string>/Users/alex/Library/Logs/sessionskills-weekly.log</string>
  <key>StandardErrorPath</key><string>/Users/alex/Library/Logs/sessionskills-weekly.log</string>
  <key>WorkingDirectory</key><string>/Users/alex/A/SessionManager</string>
</dict>
</plist>
