# SessionManager/launchd/ai.sessionskills.analyze.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>ai.sessionskills.analyze</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/alex/A/SessionManager/bin/sm-pipeline</string>
    <string>run</string>
    <string>name</string><string>analyze</string><string>write_to_memory</string>
    <string>--limit</string><string>50</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Minute</key><integer>5</integer>
  </dict>
  <key>RunAtLoad</key><false/>
  <key>StandardOutPath</key><string>/Users/alex/Library/Logs/sessionskills-analyze.log</string>
  <key>StandardErrorPath</key><string>/Users/alex/Library/Logs/sessionskills-analyze.log</string>
  <key>WorkingDirectory</key><string>/Users/alex/A/SessionManager</string>
</dict>
</plist>
